# Docker/testca/files/testca.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100

[program:ocsp]
directory=/home/ca/CA/
command=/usr/bin/python3 -u /home/ca/CA/ocsp.py 
user=ocsp
group=ca
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tsa]
directory=/home/ca/TSA
command=/usr/bin/python3 -u /home/ca/TSA/tsa_server.py
user=ocsp
group=ca
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sign]
directory=/home/ca/CA/
command=/usr/bin/python3 -u /home/ca/CA/ca.py
user=ca
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uwsgi]
directory=/var/www/acme2certifier
command=/usr/bin/uwsgi --ini /var/www/acme2certifier/acme2certifier.ini
user=www-data
group=www-data
autorestart=true
