# src/security-server/admin-service/ui-system-test/src/intTest/resources/ca-container/files/ca-xroad.conf
[supervisord]
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
priority=100

[program:ocsp]
directory=/home/ca/CA/
command=/usr/bin/python3 -u /home/ca/CA/ocsp.py
user=ocsp
group=ca
autorestart=true

[program:tsa]
directory=/home/ca/TSA
command=/usr/bin/python3 -u /home/ca/TSA/tsa_server.py
user=ocsp
group=ca
autorestart=true

[program:sign]
directory=/home/ca/CA/
command=/usr/bin/python3 -u /home/ca/CA/ca.py
user=ca
autorestart=true
