# Docker/centralserver/files/cs-xroad.conf
[supervisord]
user=root

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
stopsignal=INT
stopwaitsecs=30
autorestart=unexpected
priority=100

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
priority=100

[program:xroad-configuration-service]
command=/usr/share/xroad/bin/xroad-configuration-service
user=xroad
startsecs=5
autorestart=true
stopasgroup=true
redirect_stderr=true
priority=80

[program:xroad-center]
command=/usr/share/xroad/bin/xroad-centralserver-admin-service
user=xroad
autorestart=true

[program:xroad-center-registration-service]
command=/usr/share/xroad/bin/xroad-centralserver-registration-service
user=xroad
autorestart=true

[program:xroad-center-management-service]
command=/usr/share/xroad/bin/xroad-centralserver-management-service
user=xroad
autorestart=true

[program:xroad-signer]
command=/usr/share/xroad/bin/xroad-signer
user=xroad
autorestart=true
priority=200

[program:xroad-autologin]
command=/usr/share/xroad/autologin/xroad-autologin-retry.sh
user=xroad
autorestart=false

[program:cron]
command=/usr/sbin/cron -f
user=root
autorestart=true

[program:xroad-edc-ih]
command=/usr/share/xroad/bin/xroad-edc-identity-hub
user=root
startsecs=5
autorestart=true
priority=250

# [program:xroad-edc-connector]
# command=/usr/share/xroad/bin/xroad-edc-connector
# user=root
# startsecs=5
# autorestart=true
# priority=300
# stdout_logfile=/var/log/xroad/edc-connector.log

[program:xroad-edc-credential-service]
command=/usr/share/xroad/bin/xroad-edc-credential-service
user=root
startsecs=5
autorestart=true
priority=300
