# Docker/securityserver/files/ss-hwtoken-xroad.conf
[program:token-login]
command=/usr/share/xroad/autologin/login-inactive-token.sh
user=root
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=200
