# Docker/securityserver/files/ss-xroad.conf
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/16/main -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
stopsignal=INT
stopwaitsecs=30
autorestart=unexpected
priority=100

[program:xroad-proxy-ui-api]
command=/usr/share/xroad/bin/xroad-proxy-ui-api
user=xroad
autorestart=true

[program:xroad-signer]
command=/usr/share/xroad/bin/xroad-signer
user=xroad
autorestart=true
priority=200

[program:xroad-confclient]
command=/usr/share/xroad/bin/xroad-confclient
user=xroad
autorestart=true
priority=100

[program:xroad-monitor]
command=/usr/share/xroad/bin/xroad-monitor
user=xroad
autorestart=true

[program:xroad-proxy]
command=/usr/share/xroad/bin/xroad-proxy
user=xroad
autorestart=true
stopwaitsecs=30

[program:xroad-autologin]
command=/usr/share/xroad/autologin/xroad-autologin-retry.sh
user=xroad
autorestart=false

[program:xroad-opmonitor]
command=/usr/share/xroad/bin/xroad-opmonitor
user=xroad
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-addon-messagelog]
command=/usr/share/xroad/bin/xroad-messagelog-archiver
user=xroad
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
