# xtee6/packages/xroad/debian/xroad-confproxy.cron.d
#!/bin/sh

* * * * * xroad /usr/share/xroad/scripts/confproxy-update.sh  2>&1 >/tmp/confproxy.last;

