# xtee6/packages/xroad/debian/xroad-centralserver.postinst
#!/bin/sh

# disable confclient start
echo manual > /etc/init/xroad-confclient.override
invoke-rc.d xroad-confclient stop

exit 0
