# xtee6/packages/xroad/debian/xroad-addon-hwtokens.postinst
#!/bin/sh

invoke-rc.d xroad-signer restart

exit 0
