# xtee6/packages/xroad/debian/xroad-monitor.postinst
#!/bin/bash
invoke-rc.d xroad-monitor restart
exit 0
