# xtee6/packages/xroad/debian/trusty/xroad-addon-monitor.xroad-monitor.upstart

description	"X-Road monitor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 007
setuid sdsb
setgid sdsb
console log


pre-start script
    mkdir -p -m0755 /var/run/sdsb
    chown sdsb:sdsb /var/run/sdsb
end script

script
  set -x
  . /etc/sdsb/services/monitor-agent.conf
  exec ${JAVA_HOME}/bin/java ${MONITOR_PARAMS} ${SDSB_PARAMS} -cp ${CP}${ADDON_CP} ${ADDON_PARAMS} ee.cyber.sdsb.proxymonitoragent.Main
end script

