# xtee6/packages/xroad/redhat/SOURCES/signer/xroad-signer.service
[Unit]
Description=X-Road signer
After=network.target

[Service]
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p -m0755 /var/run/xroad
ExecStartPre=/bin/chown xroad:xroad /var/run/xroad
User=xroad
Group=xroad
ExecStart=/usr/bin/xroad-signer
Restart=on-failure

[Install]
WantedBy=multi-user.target

