# deployment/native-packages/src/xroad/ubuntu/generic/xroad-securityserver-ee.postinst
#!/bin/bash

if [ "$1" = configure ]; then
  invoke-rc.d --quiet xroad-proxy try-restart || true
  invoke-rc.d --quiet xroad-signer try-restart || true
  invoke-rc.d --quiet xroad-proxy-ui-api try-restart || true
fi

#DEBHELPER#

exit 0
