# src/packages/src/xroad/ubuntu/generic/xroad-addon-wsdlvalidator.postinst
#!/bin/bash

if [[ "$1" == 'configure' && -z "$2" ]]; then
  # set initial value for wsdl-validator-command in first install
  crudini --set /etc/xroad/conf.d/local.ini proxy-ui-api wsdl-validator-command /usr/share/xroad/wsdlvalidator/bin/wsdlvalidator_wrapper.sh
  chown root:xroad /usr/share/xroad/wsdlvalidator/bin/wsdlvalidator_wrapper.sh
  chmod 0750 /usr/share/xroad/wsdlvalidator/bin/wsdlvalidator_wrapper.sh
  invoke-rc.d --quiet xroad-proxy-ui-api try-restart || true
fi

if [ "$1" = abort-upgrade ]; then
  exit 0
fi

exit 0
