# src/packages/src/xroad/ubuntu/generic/xroad-addon-opmonitoring.postinst
#!/bin/sh

invoke-rc.d --quiet xroad-proxy try-restart || true

exit 0

