# src/packages/src/xroad/ubuntu/generic/xroad-addon-hwtokens.postinst
#!/bin/bash

invoke-rc.d --quiet xroad-signer try-restart || true

exit 0
