# src/packages/src/xroad/ubuntu/generic/xroad-centralserver.postinst
#!/bin/sh

# disable confclient start
if [ "$(lsb_release -sc)" = "trusty" ]; then
  echo manual > /etc/init/xroad-confclient.override
  service xroad-confclient stop
else [ "$(lsb_release -sc)" = "bionic" ]
  systemctl stop xroad-confclient
  systemctl disable xroad-confclient
fi

exit 0
