# src/packages/src/xroad/ubuntu/generic/xroad-addon-hwtokens.postinst
#!/bin/bash
if [ "$1" = configure ]; then
  invoke-rc.d --quiet xroad-signer try-restart || true
fi

if [ "$1" = abort-upgrade ]; then
  exit 0
fi

exit 0
