# src/packages/src/xroad/ubuntu/generic/xroad-configuration-service.service
[Unit]
Description=X-Road Configuration Service
After=network.target xroad-secret-store-local.service xroad-secret-store-remote.service
Wants=xroad-secret-store-local.service xroad-secret-store-remote.service

[Service]
User=xroad
Group=xroad
ExecStart=/usr/share/xroad/bin/xroad-configuration-service
Restart=on-failure
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
