# src/packages/src/xroad/ubuntu/generic/xroad-securityserver-is.postinst
#!/bin/bash

if [ "$1" = configure ]; then
  invoke-rc.d --quiet xroad-addon-messagelog try-restart || true
fi

#DEBHELPER#

exit 0
