# src/packages/src/xroad/ubuntu/generic/xroad-addon-hwtokens.postinst
#!/bin/bash

invoke-rc.d --quiet xroad-signer try-restart || true

if [ "$1" = abort-upgrade ]; then
  exit 0
fi

exit 0
