# src/packages/src/xroad/ubuntu/generic/xroad-base.service
[Unit]
Description=X-Road initialization

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
