# src/packages/src/xroad/ubuntu/generic/xroad-addon-messagelog.postinst
#!/bin/bash
. /usr/share/debconf/confmodule

if [[ "$1" == configure ]]; then
  RET=
  db_get xroad-addon-messagelog/enable-messagelog || RET=true

  if [ "$RET" = false ]; then
    if [ -e /usr/share/xroad/jlib/addon/proxy/messagelog.conf ]; then
      mv /usr/share/xroad/jlib/addon/proxy/messagelog.conf /usr/share/xroad/jlib/addon/proxy/messagelog.conf.disabled
    fi
  else
    if [ -e /usr/share/xroad/jlib/addon/proxy/messagelog.conf.disabled ]; then
      mv /usr/share/xroad/jlib/addon/proxy/messagelog.conf.disabled /usr/share/xroad/jlib/addon/proxy/messagelog.conf
    fi

    RET=
    db_get xroad-common/database-host || RET=""
    /usr/share/xroad/scripts/setup_messagelog_db.sh "$RET"
  fi

  db_stop

  invoke-rc.d xroad-proxy try-restart || true
fi

if [ "$1" = abort-upgrade ]; then
  exit 0
fi

exit 0
