# src/packages/src/xroad/ubuntu/generic/xroad-centralserver.postinst
#!/bin/bash

# disable confclient start
if [ "$(lsb_release -sc)" = "bionic" ]; then
  systemctl stop xroad-confclient
  systemctl disable xroad-confclient
fi

exit 0
