# src/packages/src/xroad/ubuntu/generic/xroad-center.cron.d
#!/bin/sh

* * * * * xroad curl http://127.0.0.1:8084/managementservice/gen_conf  2>&1 >/dev/null;

