# src/packages/src/xroad/ubuntu/generic/xroad-secret-store-local.service
[Unit]
Description=X-Road OpenBao Auto Unseal Service
After=network.target openbao.service
Requires=openbao.service
BindsTo=openbao.service

[Service]
Type=oneshot
User=xroad
Group=xroad
ExecStart=/usr/share/xroad/scripts/secret-store-unseal.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
