# src/packages/src/xroad/ubuntu/generic/xroad-confclient.service
[Unit]
Description=X-Road confclient
After=network.target xroad-base.service
Requires=xroad-base.service

[Service]
User=xroad
Group=xroad
ExecStart=/usr/share/xroad/bin/xroad-confclient
Restart=on-failure
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target

