# src/packages/src/xroad/ubuntu/generic/xroad-ds-catalog-service.service
[Unit]
Description=X-Road DS Catalog Service
After=network.target
Wants=postgresql.service
Wants=xroad-signer.service

[Service]
User=xroad
Group=xroad
ExecStart=/usr/share/xroad/bin/xroad-ds-catalog-service
Restart=on-failure
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
