# ansible/ss_cluster/roles/slave/files/xroad-sync.timer
[Unit]
Description=Sync X-Road configuration

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target

