# src/packages/src/xroad-jetty9/ubuntu/trusty/debian/xroad-jetty9.xroad-jetty.upstart
description     "X-Road Jetty server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 10
umask 027

setuid xroad
setgid xroad
console log
chdir /usr/share/xroad/jetty9


pre-start script
    mkdir -p -m0755 /var/run/xroad
end script

exec /usr/share/xroad/bin/xroad-jetty9

