# src/packages/xroad/debian/xroad-monitor.postinst
#!/bin/bash

chown xroad:xroad /etc/xroad/backup.d/2_xroad-monitor
chmod 0440 /etc/xroad/backup.d/2_xroad-monitor

invoke-rc.d xroad-monitor restart
exit 0
