# src/packages/xroad/debian/xroad-addon-proxymonitor.postinst
#!/bin/sh

invoke-rc.d xroad-proxy restart 
exit 0

