# Docker/conf-proxy/files/cp-xroad.conf
[program:xroad-configuration-service]
command=/usr/share/xroad/bin/xroad-configuration-service
user=xroad
startsecs=5
autorestart=true
stopasgroup=true
redirect_stderr=true
priority=80

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
priority=100

[program:xroad-signer]
command=/usr/share/xroad/bin/xroad-signer
user=xroad
autorestart=true
priority=200

[program:xroad-configuration-proxy]
command=/usr/share/xroad/bin/xroad-confproxy
user=xroad
autorestart=true
priority=220