# development/ansible/ss_cluster/roles/node/templates/postgresql-serverconf.service.j2
.include /lib/systemd/system/postgresql.service
[Service]
Environment=PGPORT={{ serverconf_db_port }}
Environment=PGDATA=/var/lib/pgsql/serverconf

