# src/central-server/management-service/int-test/src/intTest/resources/container-files/etc/supervisor/conf.d/xroad.conf
[supervisord]
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
priority=100

[program:xroad-center-management-service]
command=/usr/share/xroad/bin/xroad-centralserver-management-service
user=xroad
autorestart=true
priority=300
