# xtee6/packages/sdsb-proxy/debian/postinst
#!/bin/sh

. /usr/share/debconf/confmodule



case "$1" in
    configure)
    #chown sdsb:sdsb /etc/sdsb/serverconf.xml
    test -d /var/spool/sdsb && test -w /var/spool/sdsb || mkdir /var/spool/sdsb ; chown sdsb:sdsb /var/spool/sdsb
    test -d /var/cache/sdsb && test -w /var/cache/sdsb || mkdir /var/cache/sdsb ; chown sdsb:sdsb /var/cache/sdsb
    
    #add groups
    groups="sdsb-security-officer sdsb-registration-officer sdsb-service-administrator sdsb-system-administrator sdsb-system-auditor sdsb-users-administrator"

    RET=""
    db_get sdsb-common/username
    AUSER="$RET"
    

    for groupname in $groups
    do
      if ! getent group $groupname > /dev/null; then
        addgroup --system --quiet  $groupname
      fi
        adduser $AUSER $groupname
    done


    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

    
service sdsb-dfc restart || service sdsb-dfc start
service sdsb-proxy restart || service sdsb-proxy start
service sdsb-async restart || service sdsb-async start
service sdsb-jetty restart || service sdsb-jetty start
invoke-rc.d nginx restart

db_stop

exit 0
