# sidecar/files/slim-ss-xroad.conf
[supervisord]
user=root

[program:postgres]
command=/usr/lib/postgresql/12/bin/postgres -D /var/lib/postgresql/12/main -c config_file=/etc/postgresql/12/main/postgresql.conf
user=postgres
stopsignal=INT
stopwaitsecs=30
autostart=true
autorestart=unexpected
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-proxy-ui-api]
command=/usr/share/xroad/bin/xroad-proxy-ui-api
user=xroad
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-signer]
command=/usr/share/xroad/bin/xroad-signer
user=xroad
autorestart=true
priority=200
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-confclient]
command=/usr/share/xroad/bin/xroad-confclient
user=xroad
autorestart=true
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-proxy]
command=/usr/share/xroad/bin/xroad-proxy
user=xroad
autorestart=true
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:xroad-autologin]
command=/usr/share/xroad/autologin/xroad-autologin-retry.sh
user=xroad
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
