# src/packages/xroad/debian/trusty/xroad-autologin.xroad-autologin.upstart
description     "X-Road Automatic Token Login"
start on started xroad-signer
stop on stopped xroad-signer
task
console log
respawn
respawn limit 6 60
# 127 = PIN incorrect
normal exit 0 127

exec su -l -c "/usr/share/xroad/autologin/autologin.expect" xroad

